# deploy/apichecker-backup-user.timer
[Unit]
Description=Run API Checker backup daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
